# container/tc_docker/etc/supervisor.d/cron.ini
[program:crond]
command=/usr/sbin/crond -f
process_name=%(program_name)s
numproc=1
directory=/run
umask=022
user=root
startsecs=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=30
