# docker/tc_headphones/sv/headphones/finish
#!/bin/bash
## headphones/finish

echo "[APPLICATION] Headphones shutting down."

kill -9 $(ps ax | grep Headphones.py | awk '{print $1}')
wait $!
