# container/tc_headphones/init.d/headphones
#!/sbin/openrc-run
################################################################################
# Copyright (C) 2015-* Phillip R. Jaenke <talecaster@rootwyrm.com>
# All rights reserved
# 
# Licensed under CC-BY-NC-4.0
# See /LICENSE for details
################################################################################
# init.d/headphones

name="headphones"
supervisor="supervise-daemon"

description="Headphones is an automated music downloader for NZB and Torrent, written in Python."

: ${headphones_config:="/config/headphones/setup.cfg"}

my_address=$(hostname -i)
command="/opt/talecaster/headphones/Headphones.py"
command_args="--nolaunch --host=${my_address}"

depend() {
	need net dcron haveged
	use logger
	after firewall
}

start_pre() {
	/opt/talecaster/bin/identify.sh
	if [ -f /message ]; then
		cat /message
	fi
}
