# container/tc_prowlarr/etc/supervisor.d/indexer.ini
[program:prowlarr]
command=/opt/Prowlarr/Prowlarr -data=/talecaster/config -nobrowser
process_name=%(program_name)s
numproc=1
directory=/talecaster/config
umask=022
user=talecaster
startsecs=30
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=30
redirect_stderr=true
