# container/tc_qbittorrent/etc/supervisor.d/torrent.ini
[program:qbittorrent]
command=/usr/local/bin/qbittorrent-nox --profile=/talecaster/config
process_name=%(program_name)s
numproc=1
directory=/talecaster/config
umask=022
user=talecaster
startsecs=30
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=30
redirect_stderr=true
