# container/tc_radarr/init.d/radarr-mono
#!/sbin/openrc-run
################################################################################
# TaleCaster - https://github.com/rootwyrm/talecaster
# Copyright (C) 2015-* Phillip R. Jaenke <talecaster@rootwyrm.com> and its
# contributors. All rights reserved.
# 
# Licensed under CC-BY-NC-4.0
# See /LICENSE for details
################################################################################
# init.d/radarr

: ${datadir:="/talecaster/config"}
: ${sonarrbin:="/opt/Radarr/Radarr.exe"}
: ${command_user:="talecaster:media"}

name="Radarr"
description="Radarr - media management for movies"
supervisor="supervise-daemon"

command="/usr/local/bin/mono"
command_args="--debug ${sonarrbin} -data=${datadir} -nobrowser"
command_args_foreground=""

depend() {
	need net
	after firewall talecaster
}
