# docker/tc_nzbget/sv/nzbget/finish
#!/bin/bash

echo "[APPLICATION] nzbget shutting down."

/opt/nzbget/nzbget -c /config/nzbget.conf -Q
wait $!
