# docker/tc_couchpotato/sv/couchpotato/finish
#!/bin/bash
## headphones/finish

echo "[APPLICATION] CouchPotato shutting down."

kill -9 $(ps ax | grep CouchPotato.py | grep -v grep | awk '{print $1}')
wait $!
