# container/tc_docker/init.d/talecaster
#!/sbin/openrc-run
################################################################################
# TaleCaster - https://github.com/rootwyrm/talecaster
# Copyright (C) 2015-* Phillip R. Jaenke <talecaster@rootwrym.com> and its
# contributors, all rights reserved.
#
# Licensed under CC-BY-NC-4.0
# See /LICENSE for details
################################################################################
# tc_docker/init.d/talecaster

name="talecaster"
command="/opt/talecaster/bin/boot"
description="TaleCaster Startup Routines"

depend() {
	need net dcron
	use logger
	after firewall
	provide talecaster
}

stop_pre() {
	if [ -f /opt/talecaster/id.service ]; then
		appid=$(cat /opt/talecaster/id.service)
	fi
	if [ -f /opt/talecaster/bin/$appid ]; then
		ebegin "Running application prestart..."
		/opt/talecaster/bin/$appid
		eend $?
	fi
}

stop() {
	rm -f /firstboot
	rm -f /factory.reset
}
