# freebsd/ports/talecaster/sonarr/files/sonarr.in
#!/bin/sh
#
# Author: TaleCaster <talecaster@rootwyrm.com>
#
# $TaleCaster: master/ports/talecaster/sonarr/files/sonarr.in
#

# PROVIDE: sonarr
# REQUIRE: LOGIN auth ldap talecaster vpn
# KEYWORD: shutdown

. /etc/rc.subr
. /etc/rc.subr.talecaster

name="sonarr"
rcvar="sonarr_enable"

load_rc_config $name

: ${sonarr_enable="NO"}
: ${television_provider=""}
: ${sonarr_user:="talecaster"}
: ${sonarr_data_dir:="%%TALECASTER%%/sonarr"}

pidfile="/var/run/talecaster/sonarr.pid"
procname="%%PREFIX%%/bin/mono"
command="/usr/sbin/daemon"
command_args="-f ${procname} %%DATADIR%%/NzbDrone.exe --nobrowser --data=${sonarr_data_dir}"
start_precmd=sonarr_precmd

sonarr_precmd()
{
	export XDG_CONFIG_HOME=${sonarr_data_dir}

	if [ ! -d ${sonarr_data_dir} ]; then
		install -d -o ${sonarr_user} ${sonarr_data_dir}
	fi
}

run_rc_command "$1"

