# container/tc_nzbget/init.d/nzbget
#!/sbin/openrc-run
################################################################################
# TaleCaster - https://github.com/rootwyrm/talecaster
# Copyright (C) 2015-* Phillip R. Jaenke <talecaster@rootwrym.com> and its
# contributors, all rights reserved.
#
# Licensed under CC-BY-NC-4.0
# See /LICENSE for details
################################################################################
# init.d/nzbget

: ${cfgfile:="/talecaster/config/nzbget.conf"}
: ${command_user:="talecaster:media"}
: ${wait:=300}

name="nzbget"
description="NZBget NNTP downloader"

command="/opt/nzbget/nzbget"
command_args="-c $cfgfile -D"
#command_background="yes"
procname="nzbget"

start_stop_daemon_args="
	--wait $wait
	$start_stop_daemon_args"

depend() {
	need net
	after firewall talecaster
}
