# container/tc_frontend/etc/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx
process_name=%(program_name)s
numproc=1
directory=/etc/nginx
umask=022
user=root
startsecs=10
autorestart=unexpected
stopsignal=TERM
stopwaitsec=10
redirect_stderr=true
