# container/tc_sabnzbd/etc/supervisor.d/sabnzbd.ini
[program:sabnzbd]
command=/opt/talecaster/venv/bin/python /opt/sabnzbd/SABnzbd.py -f /talecaster/config/sabnzbd.ini --browser 0 --pidfile /talecaster/config/sabnzbd.pid
process_name=%(program_name)s
numproc=1
directory=/talecaster/config
umask=022
user=talecaster
startsecs=10
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
