# freebsd/ports/talecaster/lidarr/files/lidarr.in
#!/bin/sh
#
# Author: TaleCaster <talecaster@rootwyrm.com>
#
# $TaleCaster: master/ports/talecaster/lidarr/files/lidarr.in
#

# PROVIDE: lidarr
# REQUIRE: LOGIN auth ldap talecaster vpn
# KEYWORD: shutdown

. /etc/rc.subr
. /etc/rc.subr.talecaster

name="lidarr"
rcvar="lidarr_enable"

load_rc_config $name

: ${lidarr_enable="NO"}
: ${music_provider=""}
: ${lidarr_user:="talecaster"}
: ${lidarr_data_dir:="%%TALECASTER%%/lidarr"}

pidfile="/var/run/talecaster/lidarr.pid"
procname="%%PREFIX%%/bin/mono"
command="/usr/sbin/daemon"
command_args="-f ${procname} %%DATADIR%%/Lidarr.exe --nobrowser --data=${lidarr_data_dir}"
start_precmd=lidarr_precmd

lidarr_precmd()
{
	if [ ! -d ${lidarr_data_dir} ]; then
		install -d -o ${lidarr_user} ${lidarr_data_dir}
	fi
}

run_rc_command "$1"
