# docker/tc_docker/sv/cron/finish
#!/bin/sh

kill -SIGTERM $(cat /var/run/crond.pid)

if [ -f /var/run/crond.pid ]; then
	rm /var/run/crond.pid
fi
