# docker/tc_docker/sv/cron/run
#!/bin/sh

if [ ! -z $PERIODIC ]; then
	SYSTEM_TAB="-s $PERIODIC"
fi

if [ ! -z $USERCRON ]; then
	USER_TAB="-c $USERCRON"
fi

exec /usr/sbin/crond -l 3 $SYSTEM_TAB $USER_TAB
